# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:hermes-gateway]
command=/opt/hermes-venv/bin/hermes gateway run --profile as-tech
autostart=true
autorestart=true
startretries=5
startsecs=10
stderr_logfile=/var/log/as-tech/hermes-gateway.err.log
stdout_logfile=/var/log/as-tech/hermes-gateway.out.log
environment=HERMES_HOME="/etc/hermes",PATH="/usr/local/bin:/usr/bin:/bin:/opt/hermes-venv/bin"
priority=10

[program:n8n]
command=n8n start
user=node
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/as-tech/n8n.err.log
stdout_logfile=/var/log/as-tech/n8n.out.log
priority=20
environment=N8N_PORT="5678",N8N_HOST="0.0.0.0"
